# workdir/meta-bee-distro/recipes-apps/clock-backend/files/clock-backend.service
[Unit]
Description=Clock backend WebSocket/REST service
After=network.target data-init.service
Requires=data-init.service

[Service]
Type=simple
WorkingDirectory=/data/clock-backend
ExecStart=/usr/bin/clock-backend
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
