# workdir/meta-bee-distro/recipes-apps/clock-app/files/clock-app.service
[Unit]
Description=Service for starting the clock application.
# No dependencies required

[Service]
Type=simple
ExecStart=/usr/bin/clock-app

[Install]
WantedBy=default.target